# script/runit/nginx/run
#!/bin/sh -e
source /etc/envvars
sv check cgit >/dev/null || true
exec nginx -g 'daemon off;'
