# script/runit/cgit/run
#!/bin/sh -e
source /etc/envvars
cgit-fcgiwrap
